# rootfs/etc/supervisor.d/carbon-cache.ini
[program:carbon-cache]

autorestart             = true

command                 = /opt/graphite/bin/carbon-cache.py --debug start

stdout_logfile          = /proc/self/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile          = /proc/self/fd/2
stderr_logfile_maxbytes = 0

