# rootfs/etc/supervisor.d/carbon-cache.ini
[program:carbon-cache]
command        = /opt/graphite/bin/carbon-cache.py --debug start
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
autorestart    = true
