# rootfs/etc/supervisor.d/graphite-webapp.ini
[program:graphite-webapp]
directory      = /opt/graphite/webapp/graphite
environment    = PYTHONPATH='/opt/graphite/webapp'
command        = gunicorn graphite_wsgi:application
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
autorestart    = true

