# rootfs/etc/supervisor.d/graphite-webapp.ini
[program:graphite-webapp]

autorestart             = true

directory               = /opt/graphite/webapp/graphite
environment             = PYTHONPATH='/opt/graphite/webapp'
command                 = gunicorn graphite_wsgi:application

stdout_logfile          = /proc/self/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile          = /proc/self/fd/2
stderr_logfile_maxbytes = 0


